# tosd-server/src/tosd.service
[Unit]
Description=Tosd osd server

[Service]
ExecStart=/usr/bin/tosd-server
Restart=on-failure
RestartSec=5